# supervisord.conf
[supervisord]
nodaemon=true

[program:migration]
command=/app/binary

[program:main]
command=/app/main